# deployment/roles/optima/templates/web.service
[Unit]
Description={{ user }}

[Service]
ExecStart={{ daemon }} {{daemon_opts}}
Restart=always
User={{ user }}

[Install]
WantedBy=multi-user.target
