# deployment/components/templates/celery-optima.conf.j2
[program:celery-optima]
command=/usr/local/bin/celery -A server.webapp.tasks.celery worker -l ERROR
directory=/opt/optima
user={{ansible_user_id}}
